# service/polipo/run
#!/bin/sh
polipo proxyAddress=0.0.0.0 proxyPort=8123 socksParentProxy=localhost:9050
